# data/razer_test.service.in
[Unit]
Description=razer_test description

[Service]
Type=dbus
BusName=io.github.openrazer1
ExecStart=@bindir@/razer_test
Restart=on-abort

[Install]
Alias=dbus-io.github.openrazer1.service
